# net-check/etc/cron.d/net-check
# m h  dom mon dow   [user] command
*/3 *  *   *   *     root   /usr/bin/net-check > /dev/null 2>&1
